# vllm/local-inference/ansible/roles/recovery/templates/vllm-health-check.timer.j2
[Unit]
Description=Run vLLM Health Check Periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=1s

[Install]
WantedBy=timers.target