# vllm/local-inference/ansible/roles/backup/templates/inference-backup.service.j2
[Unit]
Description=Inference Server Backup Service
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/inference/backups/backup-inference.sh
StandardOutput=journal

[Install]
WantedBy=multi-user.target