# vllm/local-inference/ansible/roles/backup/templates/inference-backup.timer.j2
[Unit]
Description=Run Inference Server Backup Daily

[Timer]
OnCalendar=*-*-* 02:00:00
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target