# vllm/local-inference/ansible/roles/monitoring/templates/vllm-metrics.service.j2
[Unit]
Description=vLLM Metrics Collection Service
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/inference/monitoring/vllm-metrics.sh
StandardOutput=journal

[Install]
WantedBy=multi-user.target