# vllm/local-inference/ansible/roles/monitoring/templates/vllm-metrics.timer.j2
[Unit]
Description=Run vLLM Metrics Collection Service periodically

[Timer]
OnBootSec=1min
OnUnitActiveSec=30s
AccuracySec=1s

[Install]
WantedBy=timers.target